# backend/deploy/siliconfeed-docker.service.in
[Unit]
Description=SiliconFeed worker (Docker Compose)
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=${BACKEND_DIR}
# Первый запуск: из каталога backend выполните `docker compose up -d --build`.
# При каждой загрузке ОС поднимается уже собранный стек (без rebuild).
ExecStart=/usr/bin/docker compose up -d --remove-orphans
ExecStop=/usr/bin/docker compose stop
TimeoutStartSec=300

[Install]
WantedBy=multi-user.target
